# firmware_mod/config/cron/periodic/15min/record-rotate
#!/bin/sh
. /system/sdcard/scripts/common_functions.sh
/system/sdcard/controlscripts/recording stop
/system/sdcard/controlscripts/recording start

